# dockerfile-scrapy-base/supervisord.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/usr/local/apache-tomcat-8.5.13/bin/startup.sh